# servers/files/certbot/certbot-renewal.service
[Unit]
Description=Certbot Renewal

[Service]
ExecStart=/usr/bin/bash -c "cd /home/ubuntu/services/files && docker-compose run certbot"
